# GuindexProject/Guindex/service_scripts/GuindexStats.service
[Unit]
Description=Periodically calculate Guindex Stats
After=network.target

[Service]
Type=simple
Environment="PYTHONPATH=$PYTHONPATH:/var/www/"
ExecStart=/usr/bin/manage.py GuindexStats
User=root
Group=root
Restart=always
RestartSec=1
StartLimitBurst=0

[Install]
WantedBy=multi-user.target
