# GuindexProject/Guindex/service_scripts/GuindexAlerts.service
[Unit]
Description=Monitor Guindex Alerts
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/Guindex/GuindexProject
ExecStart=/usr/bin/python /var/www/Guindex/GuindexProject/manage.py GuindexAlerts
User=root
Group=root
Restart=always
RestartSec=1
StartLimitBurst=0

[Install]
WantedBy=multi-user.target
