# GuindexProject/Guindex/service_scripts/GuindexUserContributions.service
[Unit]
Description=Periodically calculate Guindex User Contributions
After=network.target

[Service]
User=root
WorkingDirectory=/var/www/Guindex/GuindexProject
ExecStart=/usr/bin/python /var/www/Guindex/GuindexProject/manage.py GuindexUserContributions

[Install]
WantedBy=multi-user.target
